# _superninja_startup.conf
[program:7860_python]
command=/bin/bash -c "source ~/.bashrc 2>/dev/null; exec python app.py"
directory=/workspace/InvoiceIQ
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/7860_python.err.log
stdout_logfile=/var/log/supervisor/7860_python.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2